# cookbooks/jan/files/app.conf
start on startup
exec /usr/local/bin/app
